# bin/dvbserver-proxy@.service
# dvbserver socket proxy to mumudvb
[Unit]
Requires=dvbserver-proxy@%i.socket
After=dvbserver-proxy@%i.socket

[Service]
Type=oneshot
ExecStart=/usr/bin/dvbserver-proxy %i
ExecStart=/bin/bash -c "/lib/systemd/systemd-socket-proxyd \
	--exit-idle-time=5 \
	0.0.0.0:$((%i+500))"
